# device-system-el9/src/main/resources/etc/profile.d/device-status-system.sh
#!/bin/bash

if /usr/bin/tty -s
then
  if test ${0:0:1} = "-"; then
    /usr/libexec/device-status/system
  fi
fi

